# parts/k8s/cloud-init/artifacts/docker-monitor.timer
[Unit]
Description=a timer that delays docker-monitor from starting too soon after boot
[Timer]
OnBootSec=30min
[Install]
WantedBy=multi-user.target
#EOF
