# parts/k8s/cloud-init/artifacts/kubelet-monitor.timer
[Unit]
Description=a timer that delays kubelet-monitor from starting too soon after boot
[Timer]
OnBootSec=30min
[Install]
WantedBy=multi-user.target